# raspberrypi.cron
# run the curl command every minute
* * * * * curl https://example.com?action=track